# cells/xcode/launchd/com.gaiaftcl.franklin.consciousness.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gaiaftcl.franklin.consciousness</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Library/Application Support/GaiaFTCL/bin/FranklinConsciousnessService</string>
    </array>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>LAUNCHD_JOB_LABEL</key>
        <string>com.gaiaftcl.franklin.consciousness</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/Library/Logs/GaiaFTCL/franklin.consciousness.log</string>
    <key>StandardErrorPath</key>
    <string>/Library/Logs/GaiaFTCL/franklin.consciousness.error.log</string>
</dict>
</plist>
