# cells/xcode/launchd/com.gaiaftcl.nats.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gaiaftcl.nats</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/nats-server</string>
        <string>-p</string>
        <string>4222</string>
        <string>-js</string>
        <string>-sd</string>
        <string>###NATS_STORE###</string>
    </array>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Library/Logs/GaiaFTCL/nats.log</string>
    <key>StandardErrorPath</key>
    <string>/Library/Logs/GaiaFTCL/nats.error.log</string>
</dict>
</plist>
